# metron-deployment/ansible/roles/sensor-stubs/templates/sensor-stubs-bro.service
[Unit]
 Description=Sends canned Bro telemetry to a Kafka topic

[Service]
 ExecStart={{ sensor_stubs_bin }}/start-bro-stub
