# metron-deployment/ansible/roles/sensor-stubs/templates/sensor-stubs-snort.service
[Unit]
 Description=Sends canned Snort telemetry to a Kafka topic

[Service]
 ExecStart={{ sensor_stubs_bin }}/start-snort-stub
